# infra/etc/systemd/system/telegram-bot.service.d/10-execstart.conf
[Service]
ExecStart=
ExecStart=/home/telegram-bot/venv/bin/python /opt/telegram-bot/bot.py
