# infra/etc/systemd/system/telegram-bot.service.d/30-envfile.conf
[Service]
EnvironmentFile=/opt/telegram-bot/.env
