# infra/etc/systemd/system/altis-heartbeat.service
[Unit]
Description=ALTIS heartbeat emit NOC SELFTEST

[Service]
Type=oneshot
ExecStart=/usr/local/bin/altis-heartbeat.sh
